# deploy/boxnet-report-jf.timer
[Unit]
Description=Run Boxnet J&F Report daily at 3:30 AM Brazilian Time
Requires=boxnet-report-jf.service

[Timer]
# Executa às 03:30 todos os dias (horário de Brasília)
OnCalendar=*-*-* 03:30:00
# Timezone brasileiro
Environment="TZ=America/Sao_Paulo"
# Se o sistema estiver desligado quando deveria executar, executa na próxima inicialização
Persistent=true
# Adiciona aleatoriedade de até 5 minutos para evitar sobrecarga
RandomizedDelaySec=300
# Precisão de 1 minuto
AccuracySec=1m

[Install]
WantedBy=timers.target